# debian-10.0-x86_64-proxmox/playbook/roles/ansible-initial-server/files/autogrowpart.service
[Unit]
Description=Automatically Grow Partition after resize by Proxmox.

[Service]
Type=simple
ExecStart=/bin/bash /usr/local/bin/auto_grow_partition.sh

[Install]
WantedBy=multi-user.target
